# droid.service
[Unit]
Description=Droid Client
# Order after the user session so /run/user/%U/{pulse,bus} exist before we
# start. Requires `loginctl enable-linger <user>` (install.sh sets this).
After=network-online.target user@%U.service
Wants=network-online.target user@%U.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/droid/droid-client.py
WorkingDirectory=%h/droid
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=PULSE_SERVER=unix:/run/user/%U/pulse/native

[Install]
WantedBy=multi-user.target
