# droid-wifi.service
[Unit]
Description=Droid WiFi Manager — AP fallback + config portal
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
# Generated by install.sh — paths computed from current install location.
ExecStart=/usr/bin/python3 __INSTALL_DIR__/wifi-manager.py
WorkingDirectory=__INSTALL_DIR__
Restart=always
RestartSec=10
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
